# deploy/systemd/photovault-client-ui.service
[Unit]
Description=photovault client UI
After=network-online.target photovault-clientd.service
Wants=network-online.target photovault-clientd.service
StartLimitIntervalSec=120
StartLimitBurst=20

[Service]
Type=simple
User=photovault
Group=photovault
WorkingDirectory=/opt/photovault/services/photovault-client-ui
EnvironmentFile=-/etc/photovault/photovault-client-ui.env
Environment=PYTHONUNBUFFERED=1
Environment=PYTHONDONTWRITEBYTECODE=1
ExecStart=/opt/photovault/.venv/bin/python -m photovault_client_ui.main
AmbientCapabilities=CAP_NET_BIND_SERVICE
CapabilityBoundingSet=CAP_NET_BIND_SERVICE
Restart=always
RestartSec=2
TimeoutStartSec=30
TimeoutStopSec=30
LogsDirectory=photovault
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectSystem=full

[Install]
WantedBy=multi-user.target
