# deploy/systemd/photovault-clientd.service
[Unit]
Description=photovault client daemon
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/photovault/services/photovault-clientd
ExecStart=/usr/bin/python -m photovault_clientd.main
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
